# deploy/sentinel-api.service
[Unit]
Description=Sentinel HTTP API — remote data access for Streamlit Cloud
After=network.target sentinel-daemon.service

[Service]
Type=simple
User=sentinel
Group=sentinel
WorkingDirectory=/opt/sentinel
Environment=PYTHONUNBUFFERED=1
Environment=SENTINEL_API_KEY=22680fb9c42a507353b71b07157c535bb886e03ee7bd7055
ExecStart=/opt/sentinel/.venv/bin/python -m src.api.server
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sentinel-api

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/sentinel
ProtectHome=false

[Install]
WantedBy=multi-user.target
