# deploy/sentinel-web.service
[Unit]
Description=Sentinel Next.js Web Frontend
After=network.target sentinel-api.service

[Service]
Type=simple
User=sentinel
Group=sentinel
WorkingDirectory=/opt/sentinel/web
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=NEXT_PUBLIC_API_URL=http://127.0.0.1:5252/api
ExecStart=/usr/bin/node /opt/sentinel/web/node_modules/.bin/next start
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sentinel-web

[Install]
WantedBy=multi-user.target
