# deploy/sentinel-api.service
[Unit]
Description=Sentinel HTTP API — remote data access for Streamlit Cloud
After=network.target sentinel-daemon.service

[Service]
Type=simple
User=sentinel
Group=sentinel
WorkingDirectory=/opt/sentinel
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/opt/sentinel/.env
ExecStart=/opt/sentinel/.venv/bin/python -m src.api.server
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sentinel-api

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/sentinel /home/sentinel/.cache
ProtectHome=false

[Install]
WantedBy=multi-user.target
